# Tests/L2Tests/pact/providerStates/providerStateService/provider-state.conf
[program:provider-state]

command=python /CT_TEST/providerStateService/start_service.py -DFOREGROUND